# symlinks/config/systemd/user/mbsync.service
[Unit]
Description=Mailbox synchronization service
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/mbsync -Va

[Install]
WantedBy=default.target
