# symlinks/config/systemd/user/vdirsyncer.timer
[Unit]
Description=Run vdirsyncer hourly

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
